# scripts/DEBIAN/postinst
#!/bin/bash

systemctl daemon-reload
systemctl enable screen-management.service
systemctl restart screen-management.service
